# deploy/production/cron/crontab
00 * * * *           /usr/local/bin/php /www/console.php subscription:prolongate                  >> /www/runtime/logs/cron/subscription_prolongate.log
